# syncd/scripts/persistent_mloop.conf
[program:persistent_mloop]
command=python3 /usr/bin/persistent_mloop.py
stdout_logfile=/tmp/mloop.out.log
stderr_logfile=/tmp/mloop.err.log
autostart=true
autorestart=false
